# config/supervisor/laravel-worker.conf
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /path/to/project/artisan queue:work redis --sleep=3 --tries=3 --timeout=90 --max-jobs=1000 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
stopwaitsecs=3600
numprocs=4
redirect_stderr=true
stdout_logfile=/path/to/project/storage/logs/worker.log
stdout_logfile_maxbytes=0

; Custom environment variables (optional)
; environment=QUEUE_CONNECTION=redis,APP_ENV=production
